# infra/oracle/systemd/travel-api.service
[Unit]
Description=Travel Share API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=travel-api
Group=travel-api
WorkingDirectory=/opt/travel-api
EnvironmentFile=/etc/travel-api/travel-api.env
ExecStart=/opt/travel-api/app
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/travel-api

[Install]
WantedBy=multi-user.target
