# scripts/travel-api.service
[Unit]
Description=Travel Planner Go API Server
After=network.target

[Service]
Type=simple
User=opc
WorkingDirectory=/home/opc/travel-api
ExecStart=/home/opc/travel-api/travel-api
Restart=on-failure
EnvironmentFile=/home/opc/travel-api/.env

[Install]
WantedBy=multi-user.target
